# systemd/caof-redis.service
[Unit]
Description=CAOF Redis — event bus and short-term memory
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/redis-server /etc/caof/redis.conf
ExecStop=/usr/bin/redis-cli shutdown
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/caof/redis
PrivateTmp=true

[Install]
WantedBy=multi-user.target
